# crontab
# add your cron jobs here
* * * * 1-5 echo 'test' >> /cron/jobs/test.log 2>&1